# ansible-k8s/roles/crio/files/crio.service.d/http-proxy.conf
[Service]
Environment="NO_PROXY=localhost"
# Environment="HTTP_PROXY=http://proxy.server:3128"
# Environment="HTTPS_PROXY=http://proxy.server:3128"